# .local/root/msung/etc/systemd/system/getty@tty1.service.d/skiplogin.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --skip-login --nonewline --noissue --autologin ski --noclear %I $TERM
