# .local/root/msung/etc/systemd/system/getty@tty1.service.d/skiplogin.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --delay 3 --noissue --init-string "\nProperty of Ted Johnson <ted@h8c.de>" --skip-login --login-options '-p -- ski' - $TERM
